# deploy/skiboerse.service
[Unit]
Description=Skiboerse Renningen - Gunicorn Django Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/skiboerse
Environment="DJANGO_DEBUG=False"
Environment="DJANGO_ALLOWED_HOSTS=skiboerse-renningen,localhost,127.0.0.1,192.168.178.65"
Environment="DB_USER=skiboerse"
Environment="DB_PASSWORD=skiboerse"
Environment="DB_HOST=localhost"
ExecStartPre=/bin/mkdir -p /var/log/gunicorn
ExecStartPre=/bin/chown pi:pi /var/log/gunicorn
ExecStart=/home/pi/skiboerse/venv/bin/gunicorn skiboerse.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile /var/log/gunicorn/error.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
